# setup/systemd/zerogex-oa-api.service
[Unit]
Description=ZeroGEX Options Analytics API Server
Documentation=https://github.com/zerogex-options/zerogex-oa
After=network-online.target zerogex-oa-ingestion.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/zerogex-oa
Environment="PATH=/home/ubuntu/zerogex-oa/venv/bin"
EnvironmentFile=/home/ubuntu/zerogex-oa/.env

# Execution
ExecStart=/home/ubuntu/zerogex-oa/venv/bin/uvicorn \
    src.api.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2 \
    --log-level info
ExecReload=/bin/kill -HUP $MAINPID

# Readiness probe: don't mark the unit Active until /api/health responds
# 200. Retries up to 30s while uvicorn binds. Sends X-API-Key when set
# in EnvironmentFile (the global Depends(api_key_auth) covers /api/health
# too) and falls back to anonymous when API_KEY is empty so dev installs
# without auth still boot. --retry-connrefused needs curl >= 7.52
# (Ubuntu 22.04 ships 7.81).
#
# `$$API_KEY` (with the $$ escape) is the shell expansion, not a systemd
# substitution — necessary so keys containing $/'/" pass through cleanly
# instead of being interpolated by systemd before reaching the shell.
ExecStartPost=/bin/sh -c 'if [ -n "$$API_KEY" ]; then \
    /usr/bin/curl --fail --silent --show-error --retry 30 --retry-delay 1 --retry-connrefused --max-time 5 \
        -H "X-API-Key: $$API_KEY" http://127.0.0.1:8000/api/health > /dev/null; \
  else \
    /usr/bin/curl --fail --silent --show-error --retry 30 --retry-delay 1 --retry-connrefused --max-time 5 \
        http://127.0.0.1:8000/api/health > /dev/null; \
  fi'

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=300

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
SendSIGKILL=yes

# Resource limits
MemoryMax=1G
CPUQuota=100%
# API wins contended cycles over background workers (ingestion=100, analytics/signals=80)
CPUWeight=400

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/zerogex-oa

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zerogex-oa-api

[Install]
WantedBy=multi-user.target
