# setup/systemd/zerogex-alert@.service
[Unit]
Description=ZeroGEX failure-alert dispatcher (instance: %i)
Documentation=https://github.com/zerogex-options/zerogex-oa
# %i = the failed unit name passed via OnFailure=zerogex-alert@%n.service
# (e.g. "zerogex-oa-normalizer-healthcheck.service").

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu

# Optional config (the leading `-` makes systemd tolerate a missing file —
# the dispatcher script itself also handles that case so a fresh install
# doesn't refuse to log).  See setup/systemd/zerogex-alert.env.example.
EnvironmentFile=-/etc/zerogex/alert.env

# Pass the failed unit's name to the dispatcher, which pulls a journal
# tail and sends to the configured backend (slack | sns | pagerduty |
# webhook | stderr).  Default backend = stderr, which lands in this
# unit's own journal — ideal for verifying the wiring before you turn
# on real alerting.
ExecStart=/home/ubuntu/zerogex-oa/setup/systemd/zerogex-alert.sh %i

# Alerting must complete fast; cap aggressively so a stuck webhook
# can't pile up alert-handler instances.
TimeoutStartSec=30

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# journalctl needs the journal directory readable; alert.env lives in
# /etc/zerogex/.  Both are read-only — we never write from here.
ReadOnlyPaths=/etc/zerogex /var/log/journal /run/log/journal

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zerogex-alert
