# setup/systemd/zerogex-oa-normalizer-refresh.timer
[Unit]
Description=Nightly trigger for ZeroGEX normalizer cache refresh (04:30 ET)

[Timer]
# Run daily at 04:30 server time (timezone is America/New_York, set
# by deploy/steps/010.setup).  Schedule rationale:
#   • 02:00 ET — db-maintain (daily, up to 30 min)
#   • 02:30 ET — db-vacuum-full (Sundays only, up to 90 min)
#   • 03:30 ET — logs-clear (daily, up to 10 min)
#   • 04:30 ET — normalizer-refresh (this unit)
# 04:30 sits safely after every other maintenance window — including
# the worst-case Sunday vacuum-full finish at 04:05 — and well before
# pre-market data flow picks up.  Jitter ±5 min to avoid lock contention
# in the unlikely case another job overruns.
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=5m
# Run on next boot if the box was off when the timer should have
# fired (don't skip nightly refreshes after planned downtime).
Persistent=true

[Install]
WantedBy=timers.target
