# setup/systemd/zerogex-oa-signals.service
[Unit]
Description=ZeroGEX Signal Engine
Documentation=https://github.com/zerogex-options/zerogex-oa
After=network-online.target zerogex-oa-analytics.service
Wants=network-online.target

# Operator hook for active alerting (e.g. PagerDuty, Slack, email).
# Wired to zerogex-alert@%n.service, which takes the failed unit name
# as %i and dispatches to whatever backend is configured in
# /etc/zerogex/alert.env (slack | sns | pagerduty | webhook | stderr).
# Without alert.env the dispatcher falls back to stderr — alerts still
# land in the alert unit's own journal so the wiring is verifiable.
OnFailure=zerogex-alert@%n.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/zerogex-oa

# Environment
Environment="PATH=/home/ubuntu/zerogex-oa/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/ubuntu/zerogex-oa/.env

# Execution
ExecStart=/home/ubuntu/zerogex-oa/venv/bin/python -m src.signals.main_engine
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=300

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
SendSIGKILL=yes

# Resource limits
MemoryMax=1G
CPUQuota=100%
CPUWeight=80

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/zerogex-oa

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zerogex-oa-signals

[Install]
WantedBy=multi-user.target
