# setup/systemd/zerogex-oa-db-vacuum-full.service
[Unit]
Description=ZeroGEX weekly full DB maintenance (VACUUM FULL + REINDEX)
Documentation=https://github.com/zerogex-options/zerogex-oa
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/home/ubuntu/zerogex-oa

# Running as root so the `sudo systemctl stop/start` calls inside
# `db-maintain-managed` are no-ops.  But systemd would otherwise set
# HOME=/root, which breaks the Makefile's `PSQL = PGPASSFILE=~/.pgpass
# psql ...` — the .pgpass file lives in the deploy user's home, not
# root's, so psql would fall through to a password prompt and exit
# non-zero on the first prune.  Pin HOME (and load .env / PATH) the
# same way the daily zerogex-oa-db-maintain unit does.
Environment="HOME=/home/ubuntu"
Environment="PATH=/home/ubuntu/zerogex-oa/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/ubuntu/zerogex-oa/.env

# `db-maintain-managed` stops the four zerogex services and runs
# `make db-maintain` (prune + VACUUM FULL + REINDEX + ANALYZE).  It
# also installs a shell EXIT trap that restarts services on exit, but
# that trap is not reliable when systemd kills the unit's cgroup on
# TimeoutStartSec — see ExecStopPost= below for the real safety net.
ExecStart=/usr/bin/make db-maintain-managed

# Unconditional restart of the four services after the maintenance
# unit transitions to inactive/failed (success, non-zero exit, OR
# TimeoutStartSec SIGTERM/SIGKILL).  Belt-and-braces for the Makefile
# trap, which does not survive a cgroup-wide SIGTERM cleanly — that is
# what left the platform dark for 8 hours on 2026-06-28.
# Idempotent: `systemctl start` on already-active units is a no-op.
ExecStopPost=/usr/bin/systemctl start zerogex-oa-ingestion zerogex-oa-analytics zerogex-oa-signals zerogex-oa-api

# VACUUM FULL + REINDEX across the maintained tables can take several
# hours on a production-sized chain.  Observed on 2026-06-28: option_chains
# alone consumed 74 min (41 min VACUUM + 33 min REINDEX); the previous
# 5400 s budget timed out partway through flow_by_contract.  4 h gives
# the full pass headroom; bump again if a future run still gets cut off.
TimeoutStartSec=14400

StandardOutput=journal
StandardError=journal
SyslogIdentifier=zerogex-oa-db-vacuum-full
